# deploy/com.eagleeye.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eagleeye.collector</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/opt/openjdk@25/bin/java</string>
        <string>--enable-native-access=ALL-UNNAMED</string>
        <string>-jar</string>
        <string>/opt/eagleeye/collector/eagleeye-collector.jar</string>
        <string>--spring.profiles.active=prod</string>
        <string>--logging.level.root=WARN</string>
        <string>--logging.level.com.eagleeye=INFO</string>
    </array>

    <!--
        Four independent collections Mon–Fri (Taipei time = UTC+8):
          13:40 (05:40 UTC): market index        — TWSE close data available ~13:30
          15:10 (07:10 UTC): institutional flow  — TWSE 三大法人 published ~15:00
          15:30 (07:30 UTC): TAIFEX OI           — 未平倉口數及契約金額 published ~15:00
          21:35 (13:35 UTC): margin transactions — TWSE 融資融券 published 20:30–21:30
    -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- 13:40 Taipei (05:40 UTC): market index -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>5</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>5</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>5</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>5</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>5</integer><key>Minute</key><integer>40</integer></dict>
        <!-- 15:10 Taipei (07:10 UTC): institutional flow -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>10</integer></dict>
        <!-- 15:30 Taipei (07:30 UTC): TAIFEX OI -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 21:35 Taipei (13:35 UTC): margin transactions -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>13</integer><key>Minute</key><integer>35</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/opt/eagleeye/logs/collector.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/eagleeye/logs/collector-error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>JAVA_HOME</key>
        <string>/usr/local/opt/openjdk@25</string>
    </dict>
</dict>
</plist>
