# deploy/com.eagleeye.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eagleeye.collector</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/opt/openjdk@25/bin/java</string>
        <string>--enable-native-access=ALL-UNNAMED</string>
        <!-- Class Data Sharing: speeds startup ~18% for this short-lived process.
             AutoCreateSharedArchive self-creates the archive on first run and
             regenerates it if the jar or JDK changes. -->
        <string>-XX:+AutoCreateSharedArchive</string>
        <string>-XX:SharedArchiveFile=/opt/eagleeye/collector/collector.jsa</string>
        <!-- Silence the cosmetic CDS module-handling notice (native-access flag
             is not recorded at dump time on JDK 25); class sharing still applies. -->
        <string>-Xlog:cds=off</string>
        <string>-jar</string>
        <string>/opt/eagleeye/collector/eagleeye-collector.jar</string>
        <string>--spring.profiles.active=prod</string>
        <string>--logging.level.root=WARN</string>
        <string>--logging.level.com.eagleeye=INFO</string>
    </array>

    <!--
        Five independent collections Mon–Fri (Taipei local time, Asia/Taipei).
        launchd StartCalendarInterval uses local time — NOT UTC.
          07:00: after-hours futures — 夜盤 session ends 05:00 on the attributed date
          15:05: market index        — TWSE FMTQIK afterTrading stats published ~15:00
          15:15: institutional flow  — TWSE 三大法人 published ~15:00
          15:30: TAIFEX OI           — 未平倉口數及契約金額 published ~15:00
          21:35: margin transactions — TWSE 融資融券 published 20:30–21:30
    -->
    <key>StartCalendarInterval</key>
    <array>
        <!-- 07:00 Taipei: after-hours futures (夜盤) -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <!-- 15:05 Taipei: market index -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>5</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>5</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>5</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>5</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>5</integer></dict>
        <!-- 15:15 Taipei: institutional flow -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>15</integer></dict>
        <!-- 15:30 Taipei: TAIFEX OI -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <!-- 21:35 Taipei: margin transactions -->
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>21</integer><key>Minute</key><integer>35</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/opt/eagleeye/logs/collector.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/eagleeye/logs/collector-error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>JAVA_HOME</key>
        <string>/usr/local/opt/openjdk@25</string>
    </dict>
</dict>
</plist>
