# deploy/com.eagleeye.collector.taifex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
    TAIFEX collector — TAIFEX OI, 未平倉口數及契約金額 published ~15:00
    launchd StartCalendarInterval uses local time (Asia/Taipei), NOT UTC.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.eagleeye.collector.taifex</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/opt/openjdk@25/bin/java</string>
        <string>--enable-native-access=ALL-UNNAMED</string>
        <!-- Class Data Sharing: speeds startup ~18% for this short-lived process. -->
        <string>-XX:+AutoCreateSharedArchive</string>
        <string>-XX:SharedArchiveFile=/opt/eagleeye/collector/collector.jsa</string>
        <string>-Xlog:cds=off</string>
        <string>-jar</string>
        <string>/opt/eagleeye/collector/eagleeye-collector.jar</string>
        <string>--spring.profiles.active=prod</string>
        <string>--logging.level.root=WARN</string>
        <string>--logging.level.com.eagleeye=INFO</string>
        <string>--collector=TAIFEX</string>
    </array>

    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>15</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/opt/eagleeye/logs/collector-taifex.log</string>
    <key>StandardErrorPath</key>
    <string>/opt/eagleeye/logs/collector-taifex-error.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>EnvironmentVariables</key>
    <dict>
        <key>JAVA_HOME</key>
        <string>/usr/local/opt/openjdk@25</string>
    </dict>
</dict>
</plist>
